# deploy/fastapi_rest_api.ini
[program:fastapi_rest_api]
chdir=/home/bewi/fastapi_rest_api
home=/home/bewi/fastapi_rest_api
pythonpath=~/fastapi_rest_api/ENV/bin/python
virtualenv=~/fastapi_rest_api/ENV+
directory=/home/bewi/fastapi_rest_api
socket=tcp://localhost:8888
command=/bin/bash -c '~/fastapi_rest_api/ENV/bin/uvicorn app:app --port 1028 --host 0.0.0.0'
numprocs=1
autostart=yes
autorestart=yes
stopsignal=INT
stopasgroup=true
killasgroup=true

